# build-rootfs/files/etc/init.d/podroid-ready
#!/sbin/openrc-run
description="Podroid VM ready marker (consumed by Android boot detector)"

depend() {
    need podroid-network
    after dropbear podroid-ratatoskr *
}

start() {
    # By the time we run, every other service has started (depend: after *).
    # OOM-protect Podroid's infrastructure processes so a hungry user app
    # (Firefox, Chromium, big npm builds, podman pulls) gets reaped first.
    # -1000 makes a process completely immune to OOM. Match by argv[0]
    # basename instead of pidfile to catch all variants (multiple Xvnc
    # screens, forked dropbear instances per SSH session, the agent's
    # listener children).
    for _name in dropbear Xvnc pulseaudio podroid-vsock-agent dnsmasq; do
        for _pid in $(pgrep -x "$_name" 2>/dev/null); do
            echo -1000 > "/proc/$_pid/oom_score_adj" 2>/dev/null
        done
    done

    _ratatoskr_container=ratatoskr
    if [ -r /etc/ratatoskr/env ]; then
        _configured=$(sed -n 's/^RATATOSKR_CONTAINER=//p' /etc/ratatoskr/env 2>/dev/null | tail -1)
        [ -n "$_configured" ] && _ratatoskr_container=$_configured
    fi
    if command -v docker >/dev/null 2>&1; then
        _i=0
        while [ "$_i" -lt 15 ]; do
            _running=$(docker inspect -f '{{.State.Running}}' "$_ratatoskr_container" 2>/dev/null || true)
            [ "$_running" = "true" ] && break
            _i=$((_i + 1))
            sleep 1
        done
        _running=$(docker inspect -f '{{.State.Running}}' "$_ratatoskr_container" 2>/dev/null || true)
        if [ "$_running" != "true" ]; then
            echo "Ratatoskr not running; continuing VM boot for diagnostics" > /dev/console
            docker logs --tail 120 "$_ratatoskr_container" > /dev/console 2>&1 || true
        fi
    fi

    # Emit the boot stages PodroidQemu.detectBootStage() expects.
    echo "Starting SSH..." > /dev/console
    echo "Almost ready..." > /dev/console
    echo "Ready!" > /dev/console
    return 0
}

stop() {
    return 0
}
